# scripts/update-site.timer
[Unit]
Description=Runs dxuuu.xyz's update script every 10 minutes

[Timer]
OnUnitActivateSec=10m
Persistent=true

[Install]
WantedBy=timers.target
