# scripts/update-cert.timer
[Unit]
Description=Runs dxuuu.xyz's update-cert script every 15 days

[Timer]
OnUnitActiveSec=15d
Persistent=true

[Install]
WantedBy=timers.target
